# php73-apache/files/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/bin/bash -c '/etc/init.d/ssh start'

[program:php-fpm]
command=/bin/bash -c 'php-fpm --nodaemonize'

[program:apache2]
command=/bin/bash -c "/usr/sbin/apache2ctl -DFOREGROUND"
